# example@.service
# Template unit for running multiple mcp-http instances.
# Copy to /etc/systemd/system/mcp-http@.service
#
# Each instance is identified by %i (the part after @).
# Config file: /etc/mcp-http/%i.env
#
# Example usage:
#   systemctl enable --now mcp-http@dav-mcp
#   systemctl enable --now mcp-http@mempalace
#   journalctl -u mcp-http@dav-mcp -f

[Unit]
Description=MCP HTTP Bridge (%i)
After=network.target

[Service]
Type=simple
User=mcp-http
EnvironmentFile=/etc/mcp-http/%i.env
ExecStart=/usr/local/bin/mcp-http
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
