# example.service
[Unit]
Description=MCP HTTP Bridge
After=network.target

[Service]
Type=simple
User=mcp-http
Group=mcp-http
WorkingDirectory=/opt/mcp-http
EnvironmentFile=-/etc/mcp-http/env
ExecStart=/opt/mcp-http/venv/bin/uvicorn server:app \
    --host 127.0.0.1 \
    --port 8765 \
    --no-access-log
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
