# example@.service
# Systemd template unit — one instance per upstream MCP server.
#
# Usage:
#   cp example@.service /etc/systemd/system/mcp-http@.service
#   systemctl daemon-reload
#
#   # Create per-instance env file:
#   cp example.env /etc/mcp-http/dav-mcp.env
#   # Edit /etc/mcp-http/dav-mcp.env: set MCP_HTTP_CMD, MCP_HTTP_PORT, etc.
#
#   systemctl enable --now mcp-http@dav-mcp
#   systemctl enable --now mcp-http@mempalace
#
#   systemctl status mcp-http@dav-mcp
#   journalctl -u mcp-http@dav-mcp -f

[Unit]
Description=MCP HTTP Bridge (%i)
After=network.target

[Service]
Type=simple
User=mcp-http
Group=mcp-http
WorkingDirectory=/opt/mcp-http
EnvironmentFile=/etc/mcp-http/%i.env
ExecStart=/opt/mcp-http/venv/bin/uvicorn server:app \
    --host 127.0.0.1 \
    --port ${MCP_HTTP_PORT} \
    --no-access-log
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
