# example.service
# Single instance. Copy to /etc/systemd/system/mcp-http.service
# and create /etc/mcp-http/default.env.

[Unit]
Description=MCP HTTP Bridge
After=network.target

[Service]
Type=simple
User=mcp-http
EnvironmentFile=/etc/mcp-http/default.env
ExecStart=/usr/local/bin/mcp-http
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
